# systemd/aur-fetch.timer
[Unit]
Description=Run AUR fetch every three days

[Timer]
OnUnitActiveSec=3d

[Install]
WantedBy=timers.target
